# build/packaging/systemd/csm.service
[Unit]
Description=CSM — Continuous Security Monitor Daemon
After=network.target

[Service]
Type=simple
ExecStart=/opt/csm/csm daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
WatchdogSec=300

[Install]
WantedBy=multi-user.target
