# build/packaging/systemd/csm.service
[Unit]
Description=CSM - Continuous Security Monitor Daemon
After=network.target

[Service]
Type=notify
NotifyAccess=main
ExecStart=/opt/csm/csm daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStartSec=120
WatchdogSec=300

StateDirectory=csm
StateDirectoryMode=0700
ConfigurationDirectory=csm
ConfigurationDirectoryMode=0750

[Install]
WantedBy=multi-user.target
